# recipes-swupdate/swupdate/swupdate/cn9130-cf-base/swu-run-postinsts.service
[Unit]
Description=SWUpdate postinst actions
Before=install-ssh-keys.service
Before=swupdate.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/swu-run-postinsts

[Install]
WantedBy=multi-user.target
